# python-graphite-web/graphite-web-build-index.service
[Unit]
Description=Graphite - Web Application - Index builder

[Service]
Type=oneshot
ExecStart=/usr/lib/graphite-web/build-index
